# tests/extreme/compose/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:xvfb]
command=Xvfb :1 -screen 0 1280x720x24
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/xvfb.err.log
stdout_logfile=/var/log/supervisor/xvfb.out.log

[program:x11vnc]
command=x11vnc -display :1 -forever -nopw -shared -rfbport 5900
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/x11vnc.err.log
stdout_logfile=/var/log/supervisor/x11vnc.out.log
environment=DISPLAY=":1"

[program:novnc]
command=websockify --web=/usr/share/novnc 6901 localhost:5900
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/novnc.err.log
stdout_logfile=/var/log/supervisor/novnc.out.log

[program:kodi]
command=bash -c "exec /usr/lib/*/kodi/kodi.bin -fs"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/kodi.err.log
stdout_logfile=/var/log/supervisor/kodi.out.log
environment=DISPLAY=":1",HOME="/root"
